# files/.config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=nm-applet
Terminal=false
Exec="/usr/bin/nm-applet"
